# deploy/supervisor/test/doaj-static.conf
[program:doaj-static]
command=/home/cloo/doaj-static/bin/gunicorn -c /home/cloo/doaj-static/src/doaj/deploy/doaj_static_gunicorn_config.py portality.app:app
environment= DOAJENV=test
user=cloo
directory=/home/cloo/doaj-static/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
