# deploy/supervisor/test/huey-main.conf
[program:huey-main]
command=/home/cloo/doaj/bin/python /home/cloo/doaj/bin/huey_consumer.py -v portality.tasks.consumer_main_queue.main_queue
environment= DOAJENV=test
user=cloo
directory=/home/cloo/doaj/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
