# deploy/supervisor/test/kafka-consumer.conf
[program:kafka-consumer]
command=/home/cloo/doaj/bin/python -m faust -A portality.events.kafka_consumer worker -l info
environment= DOAJENV=test
user=cloo
directory=/home/cloo/doaj/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
