# deploy/_legacy/supervisor/harvester/doaj-harvester.conf
[program:doaj-harvester]
command= /home/cloo/repl/harvester/doaj/bin/gunicorn -c /home/cloo/repl/harvester/doaj/src/doaj/deploy/doaj_gunicorn_config.py portality.app:app
environment = DOAJENV=harvester
user=cloo
directory=/home/cloo/repl/harvester/doaj/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
