# deploy/supervisor/production-background/huey-main.conf
[program:huey-main]
command=/home/cloo/doaj_python3/bin/python /home/cloo/doaj_python3/bin/huey_consumer.py -v portality.tasks.consumer_main_queue.main_queue
environment= DOAJENV=production
user=cloo
directory=/home/cloo/doaj_python3/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
