# deploy/supervisor/production/doaj.conf
[program:doaj]
command=/home/cloo/doaj_python3/bin/gunicorn -c /home/cloo/doaj/src/doaj_python3/deploy/doaj_gunicorn_config.py portality.app:app
environment = DOAJENV=production
user=cloo
directory=/home/cloo/doaj_python3/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
