# deploy/_legacy/supervisor/production-background/huey-main-production.conf
[program:huey-main-production]
command= /home/cloo/repl/production/doaj/bin/huey_consumer.py -v portality.tasks.consumer_main_queue.main_queue
environment= DOAJENV=production
user=cloo
directory=/home/cloo/repl/production/doaj/src/doaj
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
